# _archive/legacy_root/systemd/opsconductor-backend.service
[Unit]
Description=OpsConductor Monitor Backend (Uvicorn/FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=opsconductor
Group=opsconductor
WorkingDirectory=/home/opsconductor/opsconductor-monitor/CascadeProjects/windsurf-project
Environment="PATH=/home/opsconductor/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/opsconductor/.local/bin/uvicorn \
    app:app \
    --host 0.0.0.0 \
    --port 5000 \
    --workers 4 \
    --timeout-keep-alive 30 \
    --access-log \
    --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed

[Install]
WantedBy=multi-user.target
